# gitson.cron
# When placed in /etc/cron.d, runs PussyDNS periodically.
# m h dom mon dow user	command
*/15 *	* * *	root    cd /etc/dns && make
